# mimir-api.service.example
# /etc/systemd/system/mimir-api.service
[Unit]
Description=Mimir API Service - Multi-display content management system
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ryan
WorkingDirectory=/home/ryan/code/mimir-api/api-service
Environment=PYTHONUNBUFFERED=1
# Load environment variables (copy .env.example to .env and configure)
EnvironmentFile=-/home/ryan/code/mimir-api/api-service/.env
# Updated entry point to use new app factory
ExecStart=/home/ryan/code/mimir-api/api-service/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 5000
Restart=always
RestartSec=3

# Optional hardening
# ProtectSystem=full
# PrivateTmp=true
# NoNewPrivileges=true

[Install]
WantedBy=multi-user.target