# mimir-api/api-service/deploy/mimir-api.service
# /etc/systemd/system/mimir-api.service
[Unit]
Description=Mimir API Service - Multi-display content management system (Refactored)
Documentation=https://github.com/your-org/mimir-api
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=${REMOTE_USER}
Group=${REMOTE_USER}
WorkingDirectory=${REMOTE_PATH}
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=${REMOTE_PATH}

# Load environment variables from .env file
EnvironmentFile=-${REMOTE_PATH}/.env

# Updated entry point to use refactored app factory
ExecStart=${REMOTE_PATH}/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 5000 --workers 2
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
Restart=always
RestartSec=5

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=${REMOTE_PATH}/logs
ReadWritePaths=${REMOTE_PATH}/app.db
ReadWritePaths=/tmp
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mimir-api

[Install]
WantedBy=multi-user.target
